# templates/filebrowser.service.template
[Unit]
Description=File Browser - Web File Manager
After=network.target

[Service]
User={{FILEBROWSER_USER}}

# The server command simply needs to be pointed to the YAML config file.
# All other settings (port, database, auth method, etc.) are loaded from it.
ExecStart=/usr/local/bin/filebrowser --config /etc/filebrowser/filebrowser.yml

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target