# templates/filebrowser.service.template
[Unit]
Description=File Browser - Web File Manager
After=network.target

[Service]
User={{FILEBROWSER_USER}}
Group={{FILEBROWSER_USER}}

# The server command needs to be explicitly told where the config file is.
# This ensures it loads our proxy settings.
ExecStart=/usr/local/bin/filebrowser server --config /etc/filebrowser/.filebrowser.json

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target