# templates/filebrowser.service.template
[Unit]
Description=File Browser - Web File Manager
After=network.target

[Service]
User={{FILEBROWSER_USER}}
# The 'Group=' line has been removed to prevent the 216/GROUP error.

# DEFINITIVE FIX: The command to run the server is just 'filebrowser'
# with flags. There is no 'server' subcommand.
ExecStart=/usr/local/bin/filebrowser --config /etc/filebrowser/.filebrowser.json --database {{FILEBROWSER_DB_PATH}}

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target