# templates/filebrowser.service.template
[Unit]
Description=File Browser - Web File Manager
After=network.target

[Service]
User={{FILEBROWSER_USER}}
Group={{FILEBROWSER_USER}}

# The server command simply needs to be pointed to the pre-configured database.
# All settings (auth method, etc.) are now loaded from this file.
ExecStart=/usr/local/bin/filebrowser server --address 127.0.0.1 --port {{FILEBROWSER_PORT}} --database {{FILEBROWSER_DB_PATH}}

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target