# templates/filebrowser.service.template
[Unit]
Description=File Browser - Web File Manager
After=network.target

[Service]
User={{FILEBROWSER_USER}}

# DEFINITIVE FIX: The server command only needs the --config flag.
# All other settings, including the database path, are read from this file.
#ExecStart=/usr/local/bin/filebrowser --config /etc/filebrowser/filebrowser.yml
ExecStart=/usr/local/bin/filebrowser -c /etc/filebrowser/filebrowser.yml

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target