# templates/filebrowser.service.template
[Unit]
Description=File Browser - Web File Manager
After=network.target

[Service]
User={{FILEBROWSER_USER}}
Group={{FILEBROWSER_USER}}

# This is now correct. The service will read its configuration (auth method, etc.)
# directly from the initialized database file specified here.
ExecStart=/usr/local/bin/filebrowser server --address 127.0.0.1 --port {{FILEBROWSER_PORT}} --database {{FILEBROWSER_DB_PATH}}

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target